# build/Lede_lxc/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# 开机删除argon主题默认图片
rm /www/luci-static/argon/background/argon.jpg

exit 0
